# scripts/systemd/papers-hub-daily.timer.in
[Unit]
Description=Run Papers Hub daily refresh at 9:00 AM

[Timer]
OnCalendar=@ON_CALENDAR@
Persistent=true

[Install]
WantedBy=timers.target
